# root/etc/services.d/radarr/run
#!/usr/bin/with-contenv bash

cd /app/radarr/bin || exit

exec s6-setuidgid abc taskset -c 4-7 /app/radarr/bin/Radarr -nobrowser -data=/config
